# packages/multimedia/vdr-menuorg/unit/reloadcfg.service
[Unit]
Description=Video Disk Recorder reload registrations

[Service]
ExecStart=/usr/bin/reloadreg
Restart=never

[Install]
WantedBy=multi-user.target
