# packages/system/util-linux-ng/unit/swapfile.service
[Unit]
Description=Create and Active Swap File

[Service]
ExecStart=/usr/bin/swapfile.sh
Restart=always
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
